# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/binstagram/bin/uwsgi --ini /srv/binstagram/.config/uwsgi/uwsgi-app.ini
